# cron/crontab
0 */48 * * * /var/www/html/bin/console app:ingest-ongoing
